# backend/config/rules.json
{
  "rules": [
    {
      "rule_id": "PROC-001",
      "name": "Suspicious Reverse Shell Process",
      "event_type": "process",
      "condition": "in_list",
      "field": "process_name",
      "value": ["nc", "ncat", "netcat", "powershell.exe", "cmd.exe"],
      "severity": "high",
      "confidence": 85,
      "description": "Observed a process commonly abused for interactive shells or command execution.",
      "tactics": ["Execution", "Command and Control"],
      "techniques": ["T1059", "T1105"],
      "response_actions": ["kill_process", "generate_alert"]
    },
    {
      "rule_id": "CMD-001",
      "name": "Suspicious Command Intent",
      "event_type": "command",
      "condition": "contains",
      "field": "command_line",
      "value": [
        "mimikatz",
        "vssadmin delete shadows",
        "wmic shadowcopy delete",
        "powershell -enc",
        "powershell.exe -enc",
        "-encodedcommand",
        "certutil -urlcache",
        "bitsadmin /transfer",
        "net user",
        " /add",
        "schtasks /create",
        "reg add",
        "rundll32",
        "mshta",
        "curl http",
        "curl.exe http",
        "invoke-webrequest",
        "iex "
      ],
      "severity": "high",
      "confidence": 85,
      "description": "A command line matched a suspicious administration, persistence, download, or execution pattern.",
      "tactics": ["Execution", "Persistence", "Defense Evasion"],
      "techniques": ["T1059", "T1105", "T1547"],
      "response_actions": ["generate_alert", "investigate"]
    },
    {
      "rule_id": "AUTH-001",
      "name": "Multiple Failed Login Attempts",
      "event_type": "auth",
      "condition": "failed_login_threshold",
      "threshold": 5,
      "window_seconds": 300,
      "severity": "critical",
      "confidence": 90,
      "description": "Multiple failed login attempts observed for the same user within the configured time window.",
      "tactics": ["Credential Access"],
      "techniques": ["T1110"],
      "response_actions": ["isolate_host", "generate_alert"]
    },
    {
      "rule_id": "NET-001",
      "name": "Unauthorized Outbound Connection",
      "event_type": "network",
      "condition": "remote_ip_not_in_allowlist",
      "field": "remote_ip",
      "value": ["127.0.0.1", "::1"],
      "severity": "high",
      "confidence": 80,
      "description": "Outbound connection to a remote IP outside the baseline allowlist.",
      "tactics": ["Command and Control", "Exfiltration"],
      "techniques": ["T1071", "T1041"],
      "response_actions": ["block_ip", "generate_alert"]
    },
    {
      "rule_id": "FILE-001",
      "name": "Sensitive File Modification",
      "event_type": "file",
      "condition": "contains",
      "field": "path",
      "value": "sensitive",
      "severity": "medium",
      "confidence": 70,
      "description": "A monitored file containing the sensitive naming pattern was modified.",
      "tactics": ["Impact"],
      "techniques": ["T1485"],
      "response_actions": ["generate_alert"]
    },
    {
      "rule_id": "WEVENT-4625",
      "name": "Failed Login Attempt - Windows Event",
      "event_type": "security_event",
      "condition": "event_id_match",
      "field": "event_id",
      "value": 4625,
      "severity": "medium",
      "confidence": 85,
      "description": "A login attempt failed. May indicate brute force attack or credential stuffing.",
      "tactics": ["Credential Access", "Initial Access"],
      "techniques": ["T1110.001", "T1110.002", "T1110.003"],
      "response_actions": ["generate_alert"]
    },
    {
      "rule_id": "WEVENT-4672",
      "name": "Special Privileges Assigned - Windows Event",
      "event_type": "security_event",
      "condition": "event_id_match",
      "field": "event_id",
      "value": 4672,
      "severity": "high",
      "confidence": 90,
      "description": "Special administrative privileges were assigned. May indicate privilege escalation.",
      "tactics": ["Privilege Escalation"],
      "techniques": ["T1548.001", "T1548.002"],
      "response_actions": ["generate_alert", "monitor"]
    },
    {
      "rule_id": "WEVENT-4698",
      "name": "Scheduled Task Created - Windows Event",
      "event_type": "security_event",
      "condition": "event_id_match",
      "field": "event_id",
      "value": 4698,
      "severity": "high",
      "confidence": 85,
      "description": "A scheduled task was created. May indicate persistence mechanism.",
      "tactics": ["Persistence"],
      "techniques": ["T1053.005"],
      "response_actions": ["generate_alert", "monitor"]
    },
    {
      "rule_id": "WEVENT-4732",
      "name": "User Added to Local Group - Windows Event",
      "event_type": "security_event",
      "condition": "event_id_match",
      "field": "event_id",
      "value": 4732,
      "severity": "high",
      "confidence": 80,
      "description": "A user was added to a privileged group. May indicate lateral movement or privilege escalation.",
      "tactics": ["Privilege Escalation", "Lateral Movement"],
      "techniques": ["T1548.001", "T1570"],
      "response_actions": ["generate_alert", "investigate"]
    },
    {
      "rule_id": "WEVENT-4735",
      "name": "Local Group Modified - Windows Event",
      "event_type": "security_event",
      "condition": "event_id_match",
      "field": "event_id",
      "value": 4735,
      "severity": "high",
      "confidence": 85,
      "description": "A local group was modified. May indicate unauthorized access or privilege changes.",
      "tactics": ["Privilege Escalation"],
      "techniques": ["T1548.001"],
      "response_actions": ["generate_alert"]
    },
    {
      "rule_id": "WEVENT-4781",
      "name": "Account Created - Windows Event",
      "event_type": "security_event",
      "condition": "event_id_match",
      "field": "event_id",
      "value": 4781,
      "severity": "high",
      "confidence": 80,
      "description": "A user account was created. May indicate adversary creating backdoor accounts.",
      "tactics": ["Persistence"],
      "techniques": ["T1136.001"],
      "response_actions": ["generate_alert", "investigate"]
    },
    {
      "rule_id": "WSYS-7045",
      "name": "Service Installation - Windows Event",
      "event_type": "system_event",
      "condition": "event_id_match",
      "field": "event_id",
      "value": 7045,
      "severity": "high",
      "confidence": 85,
      "description": "A new service was installed on the system. May indicate malware or backdoor.",
      "tactics": ["Persistence", "Execution"],
      "techniques": ["T1543.003", "T1569.002"],
      "response_actions": ["generate_alert", "investigate"]
    },
    {
      "rule_id": "WSYS-7040",
      "name": "Service Startup Type Changed - Windows Event",
      "event_type": "system_event",
      "condition": "event_id_match",
      "field": "event_id",
      "value": 7040,
      "severity": "medium",
      "confidence": 75,
      "description": "A service startup type was changed. May indicate persistence modification.",
      "tactics": ["Persistence"],
      "techniques": ["T1547.001"],
      "response_actions": ["generate_alert"]
    },
    {
      "rule_id": "REG-001",
      "name": "Suspicious Registry Modification",
      "event_type": "registry_change",
      "condition": "registry_persistence_check",
      "field": "registry_path",
      "value": ["Run", "RunOnce", "Startup"],
      "severity": "high",
      "confidence": 85,
      "description": "Registry persistence key was modified. Common indicator of malware.",
      "tactics": ["Persistence"],
      "techniques": ["T1547.001", "T1547.008"],
      "response_actions": ["generate_alert", "investigate"]
    },
    {
      "rule_id": "REG-002",
      "name": "Services Registry Modified",
      "event_type": "registry_change",
      "condition": "registry_path_contains",
      "field": "registry_path",
      "value": "Services",
      "severity": "high",
      "confidence": 80,
      "description": "Services registry path was modified. May indicate malicious service installation.",
      "tactics": ["Persistence"],
      "techniques": ["T1543.003"],
      "response_actions": ["generate_alert"]
    },
    {
      "rule_id": "DNS-C2",
      "name": "Potential C2 DNS Activity",
      "event_type": "dns_query",
      "condition": "dns_c2_check",
      "field": "domain",
      "value": ["malware", "c2", "botnet", "exfil", "payload"],
      "severity": "critical",
      "confidence": 90,
      "description": "DNS query to potential command and control domain.",
      "tactics": ["Command and Control"],
      "techniques": ["T1071.004"],
      "response_actions": ["block_domain", "generate_alert", "isolate_host"]
    },
    {
      "rule_id": "DNS-TUNNEL",
      "name": "DNS Tunneling Detected",
      "event_type": "dns_query",
      "condition": "dns_entropy_high",
      "field": "domain",
      "severity": "high",
      "confidence": 85,
      "description": "High entropy in DNS queries may indicate DNS tunneling for data exfiltration.",
      "tactics": ["Exfiltration"],
      "techniques": ["T1041"],
      "response_actions": ["generate_alert", "block_domain"]
    },
    {
      "rule_id": "INJECTION-001",
      "name": "Suspicious Process Parent-Child Relationship",
      "event_type": "process_injection",
      "condition": "injection_parent_child",
      "field": "parent_process",
      "value": ["explorer.exe", "winlogon.exe", "svchost.exe"],
      "severity": "critical",
      "confidence": 95,
      "description": "Process spawned from unusual parent process. Strong indicator of injection attack.",
      "tactics": ["Defense Evasion", "Execution"],
      "techniques": ["T1055.001", "T1055.002", "T1055.004"],
      "response_actions": ["kill_process", "generate_alert", "investigate"]
    },
    {
      "rule_id": "INJECTION-002",
      "name": "Suspicious PowerShell Command",
      "event_type": "process_injection",
      "condition": "command_line_contains",
      "field": "cmdline",
      "value": ["-nop", "-enc", "-e ", "hidden", "-w hidden", "iex", "invoke-"],
      "severity": "high",
      "confidence": 90,
      "description": "PowerShell with obfuscation flags or encoded commands. Likely malicious.",
      "tactics": ["Defense Evasion", "Execution"],
      "techniques": ["T1059.001", "T1027"],
      "response_actions": ["kill_process", "generate_alert"]
    },
    {
      "rule_id": "INJECTION-003",
      "name": "Process Hollowing Indicators",
      "event_type": "process_injection",
      "condition": "hollow_process_check",
      "field": "process_name",
      "value": ["svchost.exe", "lsass.exe", "csrss.exe", "services.exe"],
      "severity": "critical",
      "confidence": 90,
      "description": "System process spawning with suspicious command line. Possible process hollowing.",
      "tactics": ["Defense Evasion"],
      "techniques": ["T1055.004"],
      "response_actions": ["kill_process", "generate_alert", "investigate"]
    },
    {
      "rule_id": "LAT-MOVE-001",
      "name": "Lateral Movement Indicator - PsExec Usage",
      "event_type": "security_event",
      "condition": "contains",
      "field": "raw_log",
      "value": "psexec",
      "severity": "critical",
      "confidence": 95,
      "description": "PsExec or similar tool used. Strong indicator of lateral movement.",
      "tactics": ["Lateral Movement"],
      "techniques": ["T1570"],
      "response_actions": ["isolate_host", "generate_alert"]
    },
    {
      "rule_id": "LAT-MOVE-002",
      "name": "Lateral Movement - Admin Share Access",
      "event_type": "network",
      "condition": "network_share_access",
      "field": "remote_path",
      "value": ["\\\\", "admin$", "c$", "ipc$"],
      "severity": "high",
      "confidence": 80,
      "description": "Access to admin shares detected. May indicate lateral movement.",
      "tactics": ["Lateral Movement"],
      "techniques": ["T1570"],
      "response_actions": ["generate_alert"]
    },
    {
      "rule_id": "EXFIL-001",
      "name": "Potential Data Exfiltration",
      "event_type": "network",
      "condition": "large_outbound_transfer",
      "field": "bytes_sent",
      "value": 1000000,
      "severity": "critical",
      "confidence": 85,
      "description": "Large outbound data transfer detected. May indicate data exfiltration.",
      "tactics": ["Exfiltration"],
      "techniques": ["T1041"],
      "response_actions": ["block_ip", "generate_alert", "investigate"]
    },
    {
      "rule_id": "EXFIL-002",
      "name": "Suspicious Cloud Sync Activity",
      "event_type": "file",
      "condition": "cloud_sync_unusual",
      "field": "process_name",
      "value": ["onedrive", "dropbox", "gdrive", "iclouddrive"],
      "severity": "high",
      "confidence": 75,
      "description": "Unusual activity in cloud sync service. May indicate data theft.",
      "tactics": ["Exfiltration"],
      "techniques": ["T1041"],
      "response_actions": ["generate_alert"]
    },
    {
      "rule_id": "RAT-001",
      "name": "Remote Access Tool Detected",
      "event_type": "process",
      "condition": "in_list",
      "field": "process_name",
      "value": ["teamviewer", "anydesk", "splashtop", "ammyy", "vnc"],
      "severity": "high",
      "confidence": 85,
      "description": "Remote access tool process detected. May indicate compromise.",
      "tactics": ["Command and Control"],
      "techniques": ["T1219"],
      "response_actions": ["kill_process", "generate_alert"]
    },
    {
      "rule_id": "RANSOMWARE-001",
      "name": "Ransomware Behavior - Mass File Encryption",
      "event_type": "file",
      "condition": "mass_file_encryption",
      "field": "file_extension",
      "value": [".locked", ".encrypted", ".crypted", ".ransomware"],
      "severity": "critical",
      "confidence": 95,
      "description": "Massive file encryption detected. Strong indicator of ransomware.",
      "tactics": ["Impact"],
      "techniques": ["T1486"],
      "response_actions": ["isolate_host", "kill_process", "generate_alert"]
    },
    {
      "rule_id": "RANSOMWARE-002",
      "name": "Ransomware Behavior - Shadow Copy Deletion",
      "event_type": "system_event",
      "condition": "contains",
      "field": "raw_log",
      "value": "vssadmin",
      "severity": "critical",
      "confidence": 90,
      "description": "Shadow copies deletion detected. Typical ransomware behavior.",
      "tactics": ["Impact"],
      "techniques": ["T1486"],
      "response_actions": ["isolate_host", "generate_alert"]
    },
    {
      "rule_id": "WEBSHELL-001",
      "name": "Potential Webshell Activity",
      "event_type": "file",
      "condition": "in_list",
      "field": "file_extension",
      "value": [".asp", ".aspx", ".jsp", ".php", ".jspx"],
      "severity": "critical",
      "confidence": 90,
      "description": "Web shell file detected in web server directory.",
      "tactics": ["Persistence"],
      "techniques": ["T1505.003"],
      "response_actions": ["delete_file", "generate_alert", "investigate"]
    },
    {
      "rule_id": "ROOTKIT-001",
      "name": "Kernel Driver Installation",
      "event_type": "system_event",
      "condition": "event_id_match",
      "field": "event_id",
      "value": 219,
      "severity": "critical",
      "confidence": 95,
      "description": "Kernel driver loaded. May indicate rootkit installation.",
      "tactics": ["Persistence"],
      "techniques": ["T1547.006"],
      "response_actions": ["isolate_host", "generate_alert"]
    },
    {
      "rule_id": "CMDCONTROL-001",
      "name": "Command and Control via WMI",
      "event_type": "process",
      "condition": "contains",
      "field": "cmdline",
      "value": "wmic",
      "severity": "high",
      "confidence": 80,
      "description": "WMI used for command execution. Common C2 technique.",
      "tactics": ["Command and Control"],
      "techniques": ["T1047"],
      "response_actions": ["kill_process", "generate_alert"]
    },
    {
      "rule_id": "PRIVESC-001",
      "name": "Privilege Escalation - UAC Bypass Attempt",
      "event_type": "process",
      "condition": "contains",
      "field": "cmdline",
      "value": ["fodhelper", "eventvwr", "slui", "dismhost"],
      "severity": "high",
      "confidence": 85,
      "description": "UAC bypass technique detected.",
      "tactics": ["Privilege Escalation"],
      "techniques": ["T1548.002"],
      "response_actions": ["kill_process", "generate_alert"]
    },
    {
      "rule_id": "PRIVESC-002",
      "name": "Privilege Escalation - Print Spooler Exploitation",
      "event_type": "system_event",
      "condition": "contains",
      "field": "raw_log",
      "value": "spoolsv",
      "severity": "critical",
      "confidence": 90,
      "description": "Print spooler activity detected. May indicate PrintNightmare or similar.",
      "tactics": ["Privilege Escalation"],
      "techniques": ["T1548.005"],
      "response_actions": ["kill_process", "generate_alert"]
    },
    {
      "rule_id": "CRED-THEFT-001",
      "name": "Credential Dumping Tool Detected",
      "event_type": "process",
      "condition": "in_list",
      "field": "process_name",
      "value": ["mimikatz", "gsecdump", "pass-the-hash"],
      "severity": "critical",
      "confidence": 95,
      "description": "Credential dumping tool process detected.",
      "tactics": ["Credential Access"],
      "techniques": ["T1003"],
      "response_actions": ["kill_process", "generate_alert", "isolate_host"]
    },
    {
      "rule_id": "CRED-THEFT-002",
      "name": "LSASS Memory Access",
      "event_type": "process",
      "condition": "lsass_memory_access",
      "field": "target_process",
      "value": "lsass.exe",
      "severity": "high",
      "confidence": 85,
      "description": "Suspicious process accessing LSASS memory. May indicate credential theft.",
      "tactics": ["Credential Access"],
      "techniques": ["T1003.001"],
      "response_actions": ["kill_process", "generate_alert"]
    }
  ]
}
